# qcom/rmtfs/rmtfs.service.in
[Unit]
Description=Qualcomm remotefs service
Requires=qrtr-ns.service
After=qrtr-ns.service

[Service]
ExecStart=RMTFS_PATH/rmtfs -r -P -s
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
